# packaging/upgs-node-panel.service
[Unit]
Description=UPGS Node Panel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/upgs-node-panel
EnvironmentFile=/opt/upgs-node-panel/.env
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=/root
Environment=PM2_HOME=/root/.pm2
ExecStart=/usr/local/bin/node server/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
