# packaging/upgs-node-panel.service
[Unit]
Description=UPGS Node Panel
After=network.target

[Service]
Type=simple
User=PANEL_USER
WorkingDirectory=INSTALL_DIR
EnvironmentFile=INSTALL_DIR/.env
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=PANEL_HOME
Environment=PM2_HOME=__PM2_HOME__
ExecStart=/usr/local/bin/node server/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
